# docker/cron/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
#PYTHONPATH=/path/to/package
*/1 * * * * root hello.sh
